# systemd/clownf_agent.service
[Unit]
Description=Clownfish Agent

[Service]
Type=simple
ExecStart=/usr/bin/clownf_agent
User=root

[Install]
WantedBy=multi-user.target
